# drupal/build/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:apache2]
command=docker-php-entrypoint apache2-foreground
autorestart=true
startretries=3
priority=10

[program:sshd]
#command=/usr/sbin/sshd -D
command=service ssh start
autorestart=true
priority=20

